# deploy/inbox-mcp-readonly.service.example
[Unit]
Description=Inbox read-only MCP HTTP gateway
After=network.target inbox-backend.service

[Service]
Type=simple
WorkingDirectory=/Users/jwalinshah/projects/inbox
EnvironmentFile=/Users/jwalinshah/projects/inbox/config/inbox.env
ExecStart=/Users/jwalinshah/projects/inbox/scripts/run_inbox_mcp_http_readonly.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
